# sheng-devauth/usr/lib/systemd/system/sheng-devauth.service.d/qtee.conf
[Unit]
Requires=qteesupplicant.service
After=qteesupplicant.service
